# systemd/.config/systemd/user/rotate_backgrounds.timer
[Unit]
Description=Automatically rotate desktop backgrounds

[Timer]
OnUnitActiveSec=20min
OnBootSec=20min

[Install]
WantedBy=timers.target
